# sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/2_message/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
set -e
declare -a options

CONFIG_PATH=/data/options.json

MQTT_Host="$(bashio::config 'MQTT_Host' 'core-mosquitto')"
MQTT_User="$(bashio::config 'MQTT_User')"
MQTT_Pass="$(bashio::config 'MQTT_Pass')"
MQTT_Topic="$(bashio::config 'MQTT_Topic')"
MQTT_Data="$(bashio::config 'MQTT_Data')"

PLANTNAME="$(bashio::config 'Plantname')"
InvSerial="$(bashio::config 'InvSerial')"

PLANTNAME2="$(bashio::config 'Plantname2' "${PLANTNAME}_Inverter2")"
InvSerial2="$(bashio::config 'InvSerial2')"

PLANTNAME3="$(bashio::config 'Plantname3' "${PLANTNAME}_Inverter3")"
InvSerial3="$(bashio::config 'InvSerial3')"

current_timestamp=$(date +'%Y-%m-%dT%H:%M:%S%z')
timestamp=$(date +'%Y-%m-%dT%H:%M:%S%z')

sleep 14

MESSAGE2="$(bashio::config 'Message2')"
if bashio::config.true 'Send_Message2'; then
            bashio::log.yellow "Sending User Message Two..."
            bashio::log.yellow "Publishing to MQTT Topic and Message"
            bashio::log.yellow /usr/bin/mosquitto_pub -h "$MQTT_Host" -u "$MQTT_User" -P "$MQTT_Pass" -t homeassistant/sbfspot_"$PLANTNAME"/sbfspot_"$InvSerial" -m "$MESSAGE2" -d
            bashio::log "MQTT Server Debug Response"
            /usr/bin/mosquitto_pub -h "$MQTT_Host" -u "$MQTT_User" -P "$MQTT_Pass" -t homeassistant/sbfspot_"$PLANTNAME"/sbfspot_"$InvSerial" -m "$MESSAGE2" -d
fi