# test-sbfspot/rootfs/etc/services.d/sbfspot/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the SBFspot service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================


# --- GENERATE CONFIG --
echo "PWD is current directory $(pwd)"
cd / 

echo "PWD is current directory $(pwd)"
/etc/sbfspot/generateConfig.sh /etc/sbfspot/SBFspot.cfg /etc/sbfspot/SBFspotUpload.cfg 


# ---- RUN ----

echo 'Starting daemon'


# cron
echo 'Starting cron in foreground'
exec /usr/sbin/crond -f
