# test-sbfspot/rootfs/etc/services.d/sbfspot/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the SBFspot service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

CONFIG_PATH=/data/options.json


# --- GENERATE CONFIG --
message=$(echo "PWD is current directory $(pwd)")
bashio::log.debug "${message:=""}"
cd / 

message=$(echo "PWD is current directory $(pwd)")
bashio::log.debug "${message:=""}"

# set -x

   if bashio::config "Bluetooth"; then
      bashio::log.info
      bashio::log.info ${__BASHIO_COLORS_BLUE} "[Bluetooth being setup]"
      bashio::log.info
      ./usr/bin/sbfspot/genBluetoothConfig.sh /usr/bin/sbfspot/SBFspot.cfg /usr/bin/sbfspot/SBFspotUpload.cfg
   fi

# set +x

# ---- Print Host BT Controller and sys info
bashio::log.info
bashio::log.info echo "$(uname -smrv)"
bashio::log.info

## Print the message the user supplied, defaults to "Hello World..."
## Get the 'message' key from the user config options.
message=$(echo "[Host Bluetooth MAC Address] $(bluetoothctl list)")
bashio::log.info
bashio::log.info "${message:="Hello World..."}"
bashio::log.info

# ---- RUN ----
message=$(echo "[Starting SBFspotUpload]")
bashio::log.info ${__BASHIO_COLORS_MAGENTA} "${message:=""}"
# cron
message=$(echo "[Starting cron in foreground]")
bashio::log.info "${message:=""}"
exec /usr/sbin/crond -f
