# test-sbfspot/rootfs/etc/services.d/sbfspot/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the SBFspot service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

CONFIG_PATH=/data/options.json


# --- GENERATE CONFIG --
echo "PWD is current directory $(pwd)"
cd / 

echo "PWD is current directory $(pwd)"


    # bashio::log.trace "${config.equals[0]}:" "$@"

    value=$(bashio::config "${"CONNECTION_TYPE"}")
    if ! bashio::var.equals "${.Connection_Type}" "${"Bluetooth"}"; then
    bashio::log.info
    bashio::log.info "[Bluetooth being setup]"
    bashio::log.info
    /usr/bin/sbfspot/genBluetoothConfig.sh /usr/bin/sbfspot/SBFspot.cfg /usr/bin/sbfspot/SBFspotUpload.cfg
        return "${__BASHIO_EXIT_NOK}"
    fi

    return "${__BASHIO_EXIT_OK}"
}


           
if bashio::config.true "$(connection_type: "Ethernet")"; then
   bashio::log.warning "Ethernet being setup"
   /usr/bin/sbfspot/genEthernetConfig.sh /usr/bin/sbfspot/SBFspot.cfg /usr/bin/sbfspot/SBFspotUpload.cfg
fi

# ---- Print Host BT Controller and sys info
bashio::log.info
bashio::log.info echo "$(uname -smrv)"
bashio::log.info

## Get the 'message' key from the user config options.
message=$(echo "[Host Bluetooth MAC Address] $(bluetoothctl list)")


## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info
bashio::log.info "${message:="Hello World..."}"
bashio::log.info

# ---- RUN ----

echo 'Starting daemon'


# cron
bashio::log.info echo '$('Starting cron in foreground')'
exec /usr/sbin/crond -f
