# haos-sbfspot-sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/init-sensorsgen/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: HAOS-SBF-MQTT-SensorsGen
# Generates Discovery MQTT messages for SBFspot sensors
# ==============================================================================

declare -a options

CONFIG_PATH=/data/options.json 

bashio::log.green "Starting SensorsGen..."