# haos-sbfspot-sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/init-sensorsgen/up
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

/etc/s6-overlay/s6-rc.d/init-sensorsgen/run